# msg-board.conf
#!upstart

description "Message Board"

start on startup
stop on shutdown

script
  export PATH="/usr/bin"
  exec sudo -u usclug node /home/usclug/code/app.js
end script
